# mkosi.extra/etc/systemd/system/debug-shell.service.d/20-exec-bash.conf
[Service]
Environment=HOME=/root
ExecStart=
ExecStart=/bin/bash --login
